# systemd/backup_daily.service
[Unit]
Description=A job to list /etc directory files

[Service]
Type=simple
ExecStart=/path/to/backup_daily.sh

[Install]
WantedBy=default.target