# systemd/backup_weekly.service
[Unit]
Description=A job to list /etc directory files

[Service]
Type=simple
ExecStart=/path/to/backup_weekly.sh

[Install]
WantedBy=default.target