# systemd/backup_weekly.service
[Unit]
Description=Weekly backup of the server

[Service]
Type=simple
ExecStart=/path/to/backup_weekly.sh

[Install]
WantedBy=default.target