# systemd/backup_weekly.timer
[Unit]
Description=Daily backup

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target